# ubus-scriptd/files/etc-ubus-daemon.init
#!/bin/sh /etc/rc.common

START=12

USE_PROCD=1
NAME=etc-ubus-daemon
PROG=/sbin/etc-ubus-daemon

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_close_instance
}

stop() {
	service_stop /sbin/etc-ubus-daemon
}

reload() {
	service_reload /sbin/etc-ubus-daemon
}
