# orange-rpcd/files/orange.init
#!/bin/sh /etc/rc.common

START=12

USE_PROCD=1
NAME=orangerpcd
PROG=/usr/bin/orangerpcd

start_service() {
	procd_open_instance
	procd_set_param respawn
	procd_set_param command $PROG -p /usr/lib/orange/api/
	procd_close_instance
}

stop_service() {
	service_stop $PROG
}

restart_service() {
	stop_service
	start_service
}

reload() {
	service_reload $PROG
}
